# deployment/survellience-app.service
[Unit]
Description=Survellience App CSE546

[Service]
Type=forking
ExecStart=/bin/bash /usr/bin/survellience-app.sh
User=ubuntu